# profile.d/alias_default.zsh
#default hack
alias -g ag='ag -i'
alias -g nkf='nkf -w -Lu -d -x'
alias -g sudo='sudo -E'
